# cerebellum/firstboot/cerebellum-growroot.service
[Unit]
Description=Cerebellum Grow Root Partition and Filesystem (once)
Wants=systemd-udev-settle.service
After=local-fs.target systemd-udev-settle.service
# Ensure root grows before firstboot and Docker bring-up
Before=cerebellum-firstboot.service docker.service
RequiresMountsFor=/

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cerebellum-growroot.sh
ExecStartPost=/bin/systemctl disable --now cerebellum-growroot.service
RemainAfterExit=no
ConditionPathExists=!/var/lib/cerebellum/growroot.done
StandardOutput=journal
StandardError=journal
TimeoutSec=120

[Install]
WantedBy=multi-user.target
