# cerebellum/firstboot/oled-statusd.service
[Unit]
Description=OLED Status Daemon (SH1106)
DefaultDependencies=no
# Bring up as early as possible once rootfs is up and journald is logging
After=local-fs.target systemd-journald.service
Before=basic.target network-pre.target multi-user.target cerebellum-firstboot.service

[Service]
Type=simple
Environment=OLED_ORIENTATION=portrait-inverted
Environment=OLED_CONTROLLER=auto
Environment=OLED_ADDR=auto
Environment=OLED_BOOT_JOURNAL=1
ExecStart=/usr/bin/python3 /opt/cerebellum/oled/oled_statusd.py
ExecStop=/usr/bin/python3 /opt/cerebellum/oled/oled_statusd.py --clear
Restart=always
RestartSec=1

[Install]
# Start during sysinit (earlier than basic.target)
WantedBy=sysinit.target
